# Configs/linux/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery notification service

[Service]
Type=oneshot
ExecStart=/home/reed/scripts/battery_notify.sh

