# Configs/linux/.config/systemd/user/battery-notify.timer
[Unit]
Description=Periodic battery check

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target

